# deploy/mysecretpve.service
# Sample systemd unit for MySecretPVE
#
# Install:
#   sudo cp deploy/mysecretpve.service /etc/systemd/system/mysecretpve.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now mysecretpve
#
# Adjust paths if you deploy to a different directory.

[Unit]
Description=MySecretPVE (Gunicorn)
After=network.target

[Service]
User=mysecretpve
Group=www-data
WorkingDirectory=/opt/mysecretpve/app
Environment=PYTHONPATH=/opt/mysecretpve/app
Environment=HTTP_DEBUG=0
RuntimeDirectory=mysecretpve
RuntimeDirectoryMode=0755
UMask=007
ExecStart=/opt/mysecretpve/venv/bin/gunicorn \
  --workers 2 \
  --timeout 30 \
  --bind unix:/run/mysecretpve/gunicorn.sock \
  webapp:app
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target

